# ansible/roles/kraken-ci/files/docker-gc.timer
[Unit]
Description=Run docker-gc.service once a day at 08:15:00 UTC (~midnight PST)
[Timer]
OnCalendar=*-*-* 08:15:00 UTC
